# university-collab-web.service
[Unit]
Description=University Collaboration Analytics Web
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/university-collaboration
Environment=WEB_PORT=8000
Environment=PATH=/var/www/university-collaboration/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/var/www/university-collaboration/.venv/bin/python web_app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
